# x11-drivers/displaylink-driver/files/99-displaylink.rules
# File based on udev rules generated by DisplayLink udev-installer.sh (6.x)

ENV{PRODUCT}!="17e9/*", GOTO="not_dl"

ACTION=="add|change", SUBSYSTEM=="usb", DRIVERS=="usb", IMPORT{builtin}="usb_id", ENV{DISPLAYLINK_DEVNAME}="$env{DEVNAME}", ENV{DISPLAYLINK_DEVICE_ID}="$env{ID_BUS}-$env{BUSNUM}-$env{DEVNUM}-$env{ID_SERIAL}"

ACTION=="add", SUBSYSTEM=="usb", DRIVERS=="usb", \
   ATTR{bInterfaceClass}=="ff", ATTR{bInterfaceProtocol}=="03", IMPORT{parent}="DISPLAYLINK*", \
   RUN+="/opt/displaylink/udev.sh add $root $devpath $env{DISPLAYLINK_DEVICE_ID} $env{DISPLAYLINK_DEVNAME}"

ACTION=="add", SUBSYSTEM=="usb", DRIVERS=="usb", \
   ATTR{bInterfaceClass}=="13", ATTR{bInterfaceProtocol}=="00", IMPORT{parent}="DISPLAYLINK*", \
   RUN+="/opt/displaylink/udev.sh add $root $devpath $env{DISPLAYLINK_DEVICE_ID} $env{DISPLAYLINK_DEVNAME}"

ACTION=="remove", RUN+="/opt/displaylink/udev.sh remove $root $env{DEVNAME}"

LABEL="not_dl"
